# deploy-webhook.service
[Unit]
Description=Deploy Webhook (fleet deploy trigger)
After=network.target

[Service]
Type=simple
# change this to your user
User=your-user
# provide the bearer token as a systemd credential instead of inline in this
# unit: store it root-only (e.g. /etc/fleet/deploy-webhook.token, chmod 600) and
# the webhook reads it from $CREDENTIALS_DIRECTORY/deploy-webhook-token. this
# keeps the token out of the world-readable unit file and out of `systemctl show`.
LoadCredential=deploy-webhook-token:/etc/fleet/deploy-webhook.token
# comma-separated app names this webhook may deploy
Environment=DEPLOY_WEBHOOK_APPS=your-app,your-app-staging
# bind loopback by default; only set a bridge-reachable host if a container must
# call it, and firewall the port if you do.
Environment=DEPLOY_WEBHOOK_HOST=127.0.0.1
# change this to your fleet directory
ExecStart=/usr/bin/node /path/to/fleet/deploy-webhook.cjs
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
# minimal hardening that does not interfere with driving deploys. heavier
# sandboxing (ProtectSystem=strict, NoNewPrivileges) is intentionally omitted
# because `fleet deploy` may escalate via sudo and drive docker/systemd.
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
